# xtools/dependency/workspaces-1.18.1/workspaces.service
[Unit]
Description=My Shell Script Service
After=docker.service

[Service]
Type=simple
ExecStart=/usr/local/src/workspaces_install.sh
User=root
Restart=no

[Install]
WantedBy=multi-user.target